# 80-naga.rules
ACTION=="add", SUBSYSTEM=="usb", ATTRS{idVendor}=="1532", SYMLINK+="naga", RUN+="/bin/systemctl --no-block start naga@%k.service"
